# ops/backup-agent.timer
[Unit]
Description=Run new-project cart backup daily

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
